# ccs.timer
[Unit]
Description=timer ccs.sh

[Timer]
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target